# launchd/com.example.slack-bridge-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.slack-bridge-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{PYTHON}}</string>
    <string>-m</string>
    <string>slack_bridge_mcp.watcher</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SLACK_BRIDGE_ENV_FILE</key>
    <string>{{ENV_FILE}}</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>{{REPO_DIR}}</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{LOG_DIR}}/slack-bridge-watcher.out.log</string>
  <key>StandardErrorPath</key>
  <string>{{LOG_DIR}}/slack-bridge-watcher.err.log</string>
</dict>
</plist>
